# crontab
# Run main.py at 3am every Saturday
0 3 * * 6 python /app/main.py >> /var/log/cron.log 2>&1
